# files/etc/uci-defaults/99a-tollgate-setup
#!/bin/sh
# TollGate first-boot configuration

# Enable or disable logging (set to 0 to disable)
ENABLE_LOGGING=1
LOGFILE="/tmp/tollgate-setup.log"

# Helper function to safely set UCI values with error handling
uci_safe_set() {
    local config="$1"
    local section="$2"
    local option="$3"
    local value="$4"
    
    # Check if the config exists
    if [ ! -f "/etc/config/$config" ]; then
        echo "Creating config: $config"
        touch "/etc/config/$config"
    fi
    
    # Check if the section exists
    if ! uci -q get "$config.$section" >/dev/null 2>&1; then
        # Section doesn't exist, try to create it
        if [[ "$section" == *"@"* ]]; then
            # For array sections like @dnsmasq[0], we need special handling
            # Extract the section type
            section_type=$(echo "$section" | cut -d'@' -f2 | cut -d'[' -f1)
            uci add "$config" "$section_type" >/dev/null 2>&1
        else
            # For named sections
            uci set "$config.$section=" >/dev/null 2>&1
        fi
    fi
    
    # Now set the option safely
    uci set "$config.$section.$option=$value" >/dev/null 2>&1
}

# Logging function
log_message() {
    if [ "$ENABLE_LOGGING" -eq 1 ]; then
        echo "$(date "+%Y-%m-%d %H:%M:%S") - $1" >> "$LOGFILE"
    fi
}

# Initialize log file if logging is enabled
if [ "$ENABLE_LOGGING" -eq 1 ]; then
    echo "TollGate Setup Log - Started at $(date)" > "$LOGFILE"
    log_message "Starting TollGate configuration script"
fi

# Safe delete function
uci_safe_delete() {
    if uci -q get "$1" >/dev/null 2>&1; then
        uci -q delete "$1"
    fi
}

# Safe add_list function
uci_safe_add_list() {
    local path="$1"
    local value="$2"
    
    # Check if config and section exist
    local config=$(echo "$path" | cut -d. -f1)
    local section=$(echo "$path" | cut -d. -f2)
    
    if [ ! -f "/etc/config/$config" ]; then
        echo "Creating config: $config"
        touch "/etc/config/$config"
    fi
    
    if ! uci -q get "$config.$section" >/dev/null 2>&1; then
        # Create section
        if [[ "$section" == *"@"* ]]; then
            section_type=$(echo "$section" | cut -d'@' -f2 | cut -d'[' -f1)
            uci add "$config" "$section_type" >/dev/null 2>&1
        else
            uci set "$config.$section=" >/dev/null 2>&1
        fi
    fi
    
    uci add_list "$path=$value" >/dev/null 2>&1
}

# 1. Firewall configuration
log_message "Configuring firewall rules"
if ! uci -q get firewall.tollgate_rules >/dev/null 2>&1; then
    uci add firewall include >/dev/null 2>&1
    uci_safe_set firewall @include[-1] path '/etc/config/firewall-tollgate'
    uci rename firewall.@include[-1]='tollgate_rules' >/dev/null 2>&1
    uci commit firewall
    log_message "Added firewall-tollgate include section"
else
    log_message "Firewall tollgate_rules already exists"
fi

# 2. Configure uhttpd listener ports and settings
log_message "Configuring uhttpd settings"

# Ensure uhttpd config exists
if [ ! -f "/etc/config/uhttpd" ]; then
    touch "/etc/config/uhttpd"
    uci add uhttpd uhttpd >/dev/null 2>&1
    uci rename uhttpd.@uhttpd[-1]='main' >/dev/null 2>&1
fi

# Set HTTP listeners
uci_safe_delete uhttpd.main.listen_http
uci_safe_add_list uhttpd.main.listen_http '0.0.0.0:8080'
uci_safe_add_list uhttpd.main.listen_http '[::]:8080'

# Set HTTPS listeners
uci_safe_delete uhttpd.main.listen_https
uci_safe_add_list uhttpd.main.listen_https '0.0.0.0:443'
uci_safe_add_list uhttpd.main.listen_https '[::]:443'

# Set other important settings
uci_safe_set uhttpd main redirect_https '0'
uci_safe_set uhttpd main home '/www'
uci_safe_set uhttpd main rfc1918_filter '1'
uci_safe_set uhttpd main max_requests '3'
uci_safe_set uhttpd main max_connections '100'
uci_safe_set uhttpd main script_timeout '60'
uci_safe_set uhttpd main network_timeout '30'
uci_safe_set uhttpd main http_keepalive '20'
uci_safe_set uhttpd main tcp_keepalive '1'

# Commit changes
uci commit uhttpd
log_message "Updated uhttpd configuration"

# 3. Configure DNS settings
log_message "Configuring DNS settings"
current_dns=$(uci -q get network.lan.dns)
if [ "$current_dns" = "127.0.0.1" ] || [ "$current_dns" = "::1" ] || [ -z "$current_dns" ]; then
    # Only change DNS if it's set to localhost or not set
    uci_safe_delete network.lan.dns
    uci_safe_add_list network.lan.dns '1.1.1.1'  # CloudFlare primary DNS
    uci_safe_add_list network.lan.dns '1.0.0.1'  # CloudFlare secondary DNS
    log_message "Set DNS servers to Cloudflare (1.1.1.1, 1.0.0.1)"
else
    log_message "Keeping existing DNS configuration: $current_dns"
fi
uci_safe_set network lan domain 'lan'

# Configure dnsmasq
log_message "Configuring dnsmasq"
uci_safe_set dhcp @dnsmasq[0] domainneeded '1'
uci_safe_set dhcp @dnsmasq[0] localise_queries '1' 
uci_safe_set dhcp @dnsmasq[0] rebind_protection '1'
uci_safe_set dhcp @dnsmasq[0] local '/lan/'
uci_safe_set dhcp @dnsmasq[0] domain 'lan'
uci_safe_set dhcp @dnsmasq[0] expandhosts '1'
uci commit network
uci commit dhcp
log_message "dnsmasq configuration completed"

# 4. Configure WiFi networks - Trail's Coffee Three-Tier Setup
log_message "Configuring Trail's Coffee three-tier WiFi networks"

# Generate a random 4-character suffix for the SSID and gateway name
RANDOM_SUFFIX=$(head /dev/urandom | tr -dc 'A-Z0-9' | head -c 4)
GATEWAY_NAME="TrailsCoffee-${RANDOM_SUFFIX}"

# SSID definitions for three-tier network
FREE_SSID_2G="TrailsCoffee-Free-2.4GHz"
FREE_SSID_5G="TrailsCoffee-Free-5GHz"
PREMIUM_SSID_2G="TrailsCoffee-Premium-2.4GHz"
PREMIUM_SSID_5G="TrailsCoffee-Premium-5GHz"
STAFF_SSID="TrailsCoffee-Staff"

log_message "Generated new SSIDs:"
log_message "  Free: ${FREE_SSID_2G} and ${FREE_SSID_5G}"
log_message "  Premium: ${PREMIUM_SSID_2G} and ${PREMIUM_SSID_5G}"
log_message "  Staff: ${STAFF_SSID}"

# Configure Free Tier Networks (2.4GHz and 5GHz) - Limited bandwidth, captive portal
if uci -q get wireless.default_radio0 >/dev/null; then
    # Configure Free 2.4GHz WiFi
    uci_safe_set wireless default_radio0 name 'trails_free_2g'
    uci_safe_set wireless default_radio0 ssid "${FREE_SSID_2G}"
    uci_safe_set wireless default_radio0 encryption 'none'
    uci_safe_set wireless default_radio0 disabled '0'
    uci_safe_set wireless default_radio0 network 'lan'  # Same network as captive portal
    log_message "Configured Free 2.4GHz WiFi: ${FREE_SSID_2G}"
else
    log_message "No 2.4GHz radio (default_radio0) found"
fi

if uci -q get wireless.default_radio1 >/dev/null; then
    # Configure Free 5GHz WiFi
    uci_safe_set wireless default_radio1 name 'trails_free_5g'
    uci_safe_set wireless default_radio1 ssid "${FREE_SSID_5G}"
    uci_safe_set wireless default_radio1 encryption 'none'
    uci_safe_set wireless default_radio1 disabled '0'
    uci_safe_set wireless default_radio1 network 'lan'  # Same network as captive portal
    log_message "Configured Free 5GHz WiFi: ${FREE_SSID_5G}"
else
    log_message "No 5GHz radio (default_radio1) found"
fi

# Configure Premium Tier Networks (additional interfaces on same radios)
# For premium networks, we'll create separate wireless interfaces
if uci -q get wireless.radio0 >/dev/null; then
    # Add premium 2.4GHz interface
    uci add wireless wifi-iface >/dev/null 2>&1
    PREMIUM_2G_IFACE="@wifi-iface[-1]"
    uci_safe_set wireless "$PREMIUM_2G_IFACE" device 'radio0'
    uci_safe_set wireless "$PREMIUM_2G_IFACE" name 'trails_premium_2g'
    uci_safe_set wireless "$PREMIUM_2G_IFACE" ssid "${PREMIUM_SSID_2G}"
    uci_safe_set wireless "$PREMIUM_2G_IFACE" encryption 'none'
    uci_safe_set wireless "$PREMIUM_2G_IFACE" disabled '0'
    uci_safe_set wireless "$PREMIUM_2G_IFACE" network 'lan'
    log_message "Added Premium 2.4GHz interface: ${PREMIUM_SSID_2G}"
fi

if uci -q get wireless.radio1 >/dev/null; then
    # Add premium 5GHz interface
    uci add wireless wifi-iface >/dev/null 2>&1
    PREMIUM_5G_IFACE="@wifi-iface[-1]"
    uci_safe_set wireless "$PREMIUM_5G_IFACE" device 'radio1'
    uci_safe_set wireless "$PREMIUM_5G_IFACE" name 'trails_premium_5g'
    uci_safe_set wireless "$PREMIUM_5G_IFACE" ssid "${PREMIUM_SSID_5G}"
    uci_safe_set wireless "$PREMIUM_5G_IFACE" encryption 'none'
    uci_safe_set wireless "$PREMIUM_5G_IFACE" disabled '0'
    uci_safe_set wireless "$PREMIUM_5G_IFACE" network 'lan'
    log_message "Added Premium 5GHz interface: ${PREMIUM_SSID_5G}"
fi

# Configure Staff Network (password-protected, separate interface)
# We'll use radio0 for staff network if available
if uci -q get wireless.radio0 >/dev/null; then
    # Add staff interface
    uci add wireless wifi-iface >/dev/null 2>&1
    STAFF_IFACE="@wifi-iface[-1]"
    uci_safe_set wireless "$STAFF_IFACE" device 'radio0'
    uci_safe_set wireless "$STAFF_IFACE" name 'trails_staff'
    uci_safe_set wireless "$STAFF_IFACE" ssid "${STAFF_SSID}"
    uci_safe_set wireless "$STAFF_IFACE" encryption 'sae-mixed'  # WPA3 with WPA2 fallback
    uci_safe_set wireless "$STAFF_IFACE" key 'TrailCoffee2024!'  # Default password - should be changed
    uci_safe_set wireless "$STAFF_IFACE" disabled '0'
    uci_safe_set wireless "$STAFF_IFACE" network 'lan'  # Same LAN network but no captive portal
    log_message "Added Staff network: ${STAFF_SSID}"
fi

# Enable wireless interfaces if they exist
if uci -q get wireless.radio0 >/dev/null; then
    uci_safe_set wireless radio0 disabled '0'
    log_message "Enabled radio0"
fi
if uci -q get wireless.radio1 >/dev/null; then
    uci_safe_set wireless radio1 disabled '0'
    log_message "Enabled radio1"
fi

# Commit wireless changes
uci commit wireless
log_message "Committed wireless configuration"

# 5. Configure NoDogSplash for Free and Premium networks
log_message "Configuring NoDogSplash for Trail's Coffee captive portal"
if ! uci -q get nodogsplash.@nodogsplash[0]; then
    # Create nodogsplash section if it doesn't exist
    uci add nodogsplash nodogsplash
    log_message "Created new nodogsplash configuration"
else
    log_message "NoDogSplash configuration already exists"
fi

# Set basic NoDogSplash configuration with Trail's Coffee branding
uci_safe_set nodogsplash @nodogsplash[0] enabled '1'
uci_safe_set nodogsplash @nodogsplash[0] gatewayname "Trail's Coffee WiFi Portal"
uci_safe_set nodogsplash @nodogsplash[0] gatewayinterface 'br-lan'

# Allow access to TollGate services
uci add_list nodogsplash.@nodogsplash[0].users_to_router='allow tcp port 2121'
uci add_list nodogsplash.@nodogsplash[0].users_to_router='allow tcp port 8080'

# Configure interface list - only Free and Premium networks get captive portal
# Staff network bypasses captive portal
if uci -q get wireless.trails_free_2g >/dev/null; then
    uci_safe_add_list nodogsplash @nodogsplash[0] gatewayinterface 'trails_free_2g'
fi
if uci -q get wireless.trails_free_5g >/dev/null; then
    uci_safe_add_list nodogsplash @nodogsplash[0] gatewayinterface 'trails_free_5g'
fi
if uci -q get wireless.trails_premium_2g >/dev/null; then
    uci_safe_add_list nodogsplash @nodogsplash[0] gatewayinterface 'trails_premium_2g'
fi
if uci -q get wireless.trails_premium_5g >/dev/null; then
    uci_safe_add_list nodogsplash @nodogsplash[0] gatewayinterface 'trails_premium_5g'
fi

uci commit nodogsplash
log_message "NoDogSplash configuration completed for Free and Premium networks"

# 6. Add TollGate port to firewall (if not in firewall-tollgate)
log_message "Checking for TollGate protocol port in firewall rules"
if ! grep -q "port 2121" /etc/config/firewall-tollgate 2>/dev/null; then
    # Add rule directly to main firewall if not in tollgate-specific file
    uci add firewall rule
    uci_safe_set firewall @rule[-1] name 'Allow-TollGate-Protocol'
    uci_safe_set firewall @rule[-1] src 'wan'
    uci_safe_set firewall @rule[-1] proto 'tcp'
    uci_safe_set firewall @rule[-1] dest_port '2121'
    uci_safe_set firewall @rule[-1] target 'ACCEPT'
    uci commit firewall
    log_message "Added firewall rule for TollGate protocol port 2121"
else
    log_message "TollGate protocol port already in firewall rules"
fi

# 7. Add first-login-setup hook to profile if it doesn't exist
log_message "Checking for first-login-setup hook in profile"
if ! grep -q "first-login-setup" /etc/profile; then
    # Append the first-login-setup code to the end of profile
    cat >> /etc/profile << 'EOF'

# TollGate first login setup
if [ ! -f /etc/first_login_done ] && [ -t 0 ] && [ -t 1 ]; then
    /usr/local/bin/first-login-setup
fi
EOF
    log_message "Added first-login-setup hook to /etc/profile"
else
    log_message "first-login-setup hook already exists in profile"
fi

# 8. Enable modem interfaces if they exist
log_message "Checking for modem interfaces to enable"
for section in $(uci show network | cut -d. -f2); do
    if [[ "$section" == "modem_"* ]]; then
        if uci -q get "network.$section" >/dev/null 2>&1; then
            uci_safe_set network "$section" disabled '0'
            log_message "Enabled modem interface: $section"
        fi
    fi
done
uci commit network
log_message "Committed network configuration changes"

# Final log message
log_message "TollGate setup completed successfully"

# Ensure script exits successfully
exit 0